# .devcontainer/Dockerfiles/redmine/supervisord.conf
[supervisord]
nodaemon=true

[program:redmine]
directory=/usr/src/redmine
command=/docker-entrypoint.sh rails server -b 0.0.0.0
autorestart=false
